# mr3020/files/etc/crontabs/root
0 2 * * * /radio_router/cgi-bin/sh/get_sunrise_sunset.sh
*/1 * * * *  /radio_router/cgi-bin/sh/cron_play.sh > /dev/null 2>&1
