# mr3020/files/etc/crontabs/root
*/1 * * * *  /radio_router/cgi-bin/sh/cron_play.sh > /dev/null 2>&1
0-59 22-7 * * * /radio_router/cgi-bin/sh/alarm_stop.sh > /dev/null 2>&1
0 8 * * * /radio_router/cgi-bin/sh/alarm_play.sh > /dev/null 2>&1
